# testing/tested.socket
[Socket]
ListenStream=/run/tested.sk

[Install]
WantedBy=sockets.target